# oracle.service
[Unit]
Description=Symphony Oracle
After=network.target

[Service]
Environment="USER=YOUR USERNAME"
Type=simple
User=${USER}
WorkingDirectory=/home/${USER}/symphony-oracle-voter
ExecStart=/home/${USER}/symphony-oracle-voter/venv/bin/python3 -u /home/${USER}/symphony-oracle-voter/main.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

Environment="SYMPHONYD_PATH=/home/${USER}/symphony/build"

# Required configuration values

[Install]
WantedBy=multi-user.target
